# target/device/valka/target_skeleton/etc/rc.d/gui
#!/bin/sh
. /etc/rc.subr

firsttime() {
	if [ -e /var/db/gui.done ]; then	
		return 0
	fi
	echo -n " * Running fc-cache: "
	/usr/bin/fc-cache -v > /dev/null
	echo "Ok"

	echo -n " * Querying PANGO modules: "
	/usr/bin/pango-querymodules > /etc/pango/pango.modules
	echo "Ok"

	echo -n " * Querying GTK Modules: "
	/usr/bin/gtk-query-immodules-2.0 > /etc/gtk-2.0/gtk.immodules
	echo "Ok"

	echo -n " * Querying Pixbuf loaders: "
	/usr/bin/gdk-pixbuf-query-loaders > /etc/gtk-2.0/gdk-pixbuf.loaders
	echo "Ok"

	touch /var/db/gui.done
}

start() {
	firsttime

	echo -n " * Starting XOrg: "
	${xorg_program} ${xorg_flags} > ${xorg_log} 2>&1 &
	echo "Ok"
}

stop() {
	return 0
}

rc_run_command "$1" "gui"
