# board/common/overlay/etc/crontabs/root
14 3 * * * /usr/sbin/logrotate /etc/logrotate.conf

