# target/device/valka/target_skeleton/etc/crontab/root
#min	hour	mday	month	wday	command
1	3	*	*	*	echo "Hello" >> /var/crondtest
