# package/rpcbind/rpcbind.service
[Unit]
Description=RPC bind service
Requires=rpcbind.socket

[Service]
Type=forking
ExecStart=/usr/sbin/rpcbind

[Install]
Also=rpcbind.socket
