# target/device/valka/target_skeleton/etc/inittab
::sysinit:/etc/rc

ttyS0::respawn:-/bin/bash
# Uncomment this to run a getty on the first serial port.
#ttyS0::respawn:/sbin/getty -L ttyS0 115200 vt100

# Run a script on shutdown.
::shutdown:/etc/rc.shutdown
