# target/device/valka/target_skeleton/etc/rc.local
#!/bin/ash
# Local startup file, executed after system is up and running

