# package/boinc/boinc-client.service
[Unit]
Description=BOINC client
After=network.target

[Service]
ExecStart=/usr/bin/boinc_client
Restart=always

[Install]
WantedBy=multi-user.target
