# package/rpi-userland/vcfiled.service
[Unit]
Description=VideoCore file server daemon

[Service]
ExecStart=/usr/sbin/vcfiled --foreground

[Install]
WantedBy=multi-user.target
